# j5/crontab
* * * * * curl localhost:3000 >> /Users/dany/projects/arduino/j5/data.txt
